# test.conf
description "iot server"
author "RvP"

start on syno.share.ready
stop on runlevel [06]

setuid http

exec node /iot/index.js
